# authresponder/systemd/auth-responder.service
[Unit]
Description=Nginx Auth Responder

[Service]
ExecStart=/usr/bin/auth-responder
User=auth-responder
Group=nginx
